# services/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = supervisord_mcp.conf

[program:php-fpm]
command=/usr/local/sbin/php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
priority=4
startsecs=5
command=nginx -g "daemon off;"
directory=/app
autostart=true
redirect_stderr=True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nextcloud-init]
command=/init_nextcloud.sh
user=%(ENV_NEXTCLOUD_USER)s
autostart=true
autorestart=false
startretries=3
priority=16
depends_on=php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:postgresql]
command=bash -c "PG_VERSION=$(pg_config --version | grep -oE '[0-9]+' | head -1); PID_FILE=/var/lib/postgresql/$PG_VERSION/main/postmaster.pid; if [ -f $PID_FILE ] && ! ps -p $(head -1 $PID_FILE 2>/dev/null) >/dev/null 2>&1; then rm -f $PID_FILE; fi; exec /usr/lib/postgresql/$PG_VERSION/bin/postgres -D /var/lib/postgresql/$PG_VERSION/main"
user=%(ENV_POSTGRES_USER)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=5
stopasgroup=true
killasgroup=true

[program:init_db]
command=/init_db.sh
user=%(ENV_POSTGRES_USER)s
autostart=true
autorestart=false
startretries=3
priority=6
depends_on=postgresql
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10

[program:mattermost]
command=bash -c "until [ -f /tmp/db_initialized ]; do echo 'Waiting for database to be initialized...'; sleep 5; done && /opt/mattermost/bin/mattermost"
user=%(ENV_MATTERMOST_USER)s
environment=MM_CONFIG=/opt/mattermost/config/config.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=20
priority=20
depends_on=init_db

[program:mattermost-init]
command=/init_mattermost.sh
autostart=true
autorestart=false
startretries=3
priority=21
depends_on=mattermost
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vnc]
command=/vnc-setup.sh
autostart=true
autorestart=true
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10
startretries=3
environment=USER=%(ENV_VNC_USER)s,HOME=/%(ENV_VNC_USER)s,DISPLAY=:1

[program:novnc]
command=websockify --web=/usr/share/novnc 6080 localhost:5901
autostart=true
autorestart=true
priority=26
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
depends_on=vnc

[program:drbench-setup]
command=/init_drbench.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:filebrowser]
command=/init_filebrowser.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=32
autostart=true
autorestart=true
startretries=3

[program:healthcheck]
command=/init_healthcheck.sh
autostart=true
autorestart=false
priority=100
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Email Server Components
# First run mail-init to set up configurations
[program:mail-init]
command=/init_mail.sh
autostart=true
autorestart=false
startretries=3
priority=35
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Then start Postfix after mail-init completes
[program:postfix]
command=/usr/sbin/postfix start-fg
autostart=true
autorestart=true
priority=36
depends_on=mail-init
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Finally start Dovecot after mail-init completes
[program:dovecot]
command=bash -c "until [ -f /tmp/mail_initialized ]; do echo 'Waiting for mail init...'; sleep 2; done && exec /usr/sbin/dovecot -F"
autostart=true
autorestart=true
priority=37
depends_on=mail-init
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
